# arch-install/dotconfig/systemd/user/clear_dotcache.timer
[Unit]
Description=Timer to clear the ~/.cache/ folder

[Timer]
OnBootSec=2min
OnUnitActiveSec=4d
Persistent=true

[Install]
WantedBy=default.target
